# dot_config/systemd/user/power_notify.service
[Unit]
Description=PowerNotify

[Service]
ExecStart=/usr/local/bin/power-manager

[Install]
WantedBy=default.target