# dot_config/systemd/exact_user/cbatticon.service.tmpl
{{ if eq .pc_loc "laptop" -}}
[Unit]
Description=PowerNotify

[Service]
ExecStart=/bin/cbatticon -i symbolic

[Install]
WantedBy=default.target
{{ end }}