# dot_config/systemd/exact_user/offlineimap.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/offlineimap -u quiet