# dot_config/systemd/exact_user/greenclip.service.d/override.conf
[Service]
RestartSec=60

[Unit]
StartLimitIntervalSec=0