# dot_config/systemd/exact_user/emacs.service.d/override.conf
[Service]
Environment="TMPDIR=/tmp/%u"
ExecStop=
ExecStop=/usr/bin/emacsclient --eval "(let (kill-emacs-hook)(kill-emacs))"
