# dot_config/systemd/exact_user/redshift.service.d/override.conf
[Service]
StartLimitIntervalSec=0
RestartSec=60